# conf/crontab
* * * * * /opt/dynname/bin/reload-config
